# k8s-wttt/rootfs/etc/services.d/syslog/run
#!/usr/bin/execlineb -P
syslogd -n -S -O -
